# vps/second-brain-hub/deploy/com.rbedu.mrluc-dashboard-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rbedu.mrluc-dashboard-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>REPO_PATH/scripts/watch_dashboard.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>VAULT_PATH</key>
    <string>VAULT_PATH</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/mrluc-dashboard-watch.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mrluc-dashboard-watch.err</string>
</dict>
</plist>
