# vps/second-brain-hub/deploy/crontab
# Second Brain — Europe/Prague (supercronic)
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/bin
VAULT_PATH=/data/mrluc
DASHBOARD_JSON=/data/mrluc/00-System/dashboard-data.json
LEGACY_TASKS=/data/mrluc/00-System/dashboard-tasks-source.json
CALENDAR_DAYS_AHEAD=2
TZ=Europe/Prague

# Po-Pa: triage 7:00, 14:00, 20:00
0 7,14,20 * * 1-5 /usr/local/bin/python3 /app/cron/triage_run.py >> /var/log/second-brain/triage.log 2>&1
# Po-Pa: rebuild +5 min
5 7,14,20 * * 1-5 /usr/local/bin/python3 /app/cron/build_dashboard.py >> /var/log/second-brain/build.log 2>&1
# Po-Pa: edu news 7:10
10 7 * * 1-5 /usr/local/bin/python3 /app/cron/edu_news_refresh.py >> /var/log/second-brain/edu.log 2>&1

# So-Ne: jen ranni beh
0 7 * * 0,6 /usr/local/bin/python3 /app/cron/triage_run.py >> /var/log/second-brain/triage.log 2>&1
5 7 * * 0,6 /usr/local/bin/python3 /app/cron/build_dashboard.py >> /var/log/second-brain/build.log 2>&1
10 7 * * 0,6 /usr/local/bin/python3 /app/cron/edu_news_refresh.py >> /var/log/second-brain/edu.log 2>&1
