# vps/second-brain-hub/deploy/crontab
# Second Brain v2 — Europe/Prague (supercronic)
# Stateless: vault I/O via Google Drive API (lib/drive_io.py).
# VAULT_DRIVE_ID + GOOGLE_DRIVE_OAUTH_JSON + GOOGLE_DRIVE_SA_JSON injected by Coolify.
#
# V2 architecture (since F6/F7):
# - File-per-task in 02-PROJEKTY/<slug>/tasks/*.md (not H3 in hub anymore)
# - Bases dashboard reads frontmatter live (no JSON build needed for Obsidian view)
# - Lifecycle scripts: done_from_checkboxes, waiting_to_asap, overdue_flag, archive_done, recurring
# - Legacy build_dashboard.py + sync_tasks_from_projekty.py kept only for HTML dashboard backwards-compat
#   (deletion in F9 — keep until F8 agent-context.json is in production)
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/bin
CALENDAR_DAYS_AHEAD=2
TZ=Europe/Prague

# === V2 LIFECYCLE (every 2 hours, staggered +0–+6 min) ===
# :00 — done from checkboxes (auto-flip status: Done if all [x])
0 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_done_from_checkboxes.py >> /var/log/second-brain/lifecycle.log 2>&1
# :01 — waiting → asap when waitUntil expired (clears waitUntil)
1 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_waiting_to_asap.py >> /var/log/second-brain/lifecycle.log 2>&1
# :02 — Waiting without waitUntil → default +3 days (reactivation date)
2 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_waiting_default_waituntil.py >> /var/log/second-brain/lifecycle.log 2>&1
# :03 — clear stray waitUntil on non-Waiting tasks (manual status edits)
3 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_waituntil_hygiene.py >> /var/log/second-brain/lifecycle.log 2>&1
# :04 — overdue flag (append OVERDUE log line in body)
4 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_overdue_flag.py >> /var/log/second-brain/lifecycle.log 2>&1
# :05 — archive Done tasks (move .md → 07-ARCHIV/tasks-done/<slug>/)
5 */2 * * * /usr/local/bin/python3 /app/cron/archive_done_tasks.py >> /var/log/second-brain/lifecycle.log 2>&1
# :06 — recurring tasks rotation (Done → archive + create next instance, F7.2)
6 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_recurring.py >> /var/log/second-brain/lifecycle.log 2>&1
# :07 — hub charter ## Stav (auto) blocks + staleness flags
7 */2 * * * /usr/local/bin/python3 /app/cron/lifecycle_hub_state.py >> /var/log/second-brain/lifecycle.log 2>&1

# === TRIAGE (LLM via cursor-agent CLI; manual fallback if CLI missing) ===
# Po-Pa: 7:00, 14:00, 20:00
0 7,14,20 * * 1-5 /usr/local/bin/python3 /app/cron/triage_llm_run.py >> /var/log/second-brain/triage.log 2>&1
# So-Ne: jen ranní běh
0 7 * * 0,6 /usr/local/bin/python3 /app/cron/triage_llm_run.py >> /var/log/second-brain/triage.log 2>&1
# Fallback inventory (Po 6:55 — log only, no proposals)
55 6 * * 1-5 /usr/local/bin/python3 /app/cron/inbox_inventory.py >> /var/log/second-brain/triage.log 2>&1

# === EDU NEWS (recurring extra) ===
# Po-Pa: edu news 7:10 (top kandidáti tématu pro OPS2 recurring task)
10 7 * * 1-5 /usr/local/bin/python3 /app/cron/lifecycle_extra_edu_news.py >> /var/log/second-brain/edu.log 2>&1
# So-Ne: jen ranní běh
10 7 * * 0,6 /usr/local/bin/python3 /app/cron/lifecycle_extra_edu_news.py >> /var/log/second-brain/edu.log 2>&1

# === WEEKLY / RETRO (Sunday evening) ===
0 20 * * 0 /usr/local/bin/python3 /app/cron/weekly_summary_draft.py >> /var/log/second-brain/weekly.log 2>&1
10 20 * * 0 /usr/local/bin/python3 /app/cron/retro_draft.py >> /var/log/second-brain/retro.log 2>&1

# === AGENT CONTEXT (F8) ===
# Every 15 min during workday — refresh agent-context.json
*/15 7-22 * * * /usr/local/bin/python3 /app/cron/build_agent_context.py >> /var/log/second-brain/agent-context.log 2>&1

# === ASAP BACKFILL (work hours) ===
# Hourly 10:00–23:00 + 00:00–02:00 — if open ASAP < 3, promote top Next by today_score
0 10-23,0-2 * * * /usr/local/bin/python3 /app/cron/lifecycle_asap_backfill.py >> /var/log/second-brain/lifecycle.log 2>&1

# === LEGACY (deprecated, will be removed in F9) ===
# Po-Pa: rebuild +5 min (po triage) — pouze pro stary HTML dashboard
# 5 7,14,20 * * 1-5 /usr/local/bin/python3 /app/cron/build_dashboard.py >> /var/log/second-brain/build.log 2>&1
# 35 8-21 * * * /usr/local/bin/python3 /app/cron/build_dashboard.py >> /var/log/second-brain/build-hourly.log 2>&1
